# systemd/minidlna.service.in
[Unit]
Description=minidlna server
After=network.target

[Service]
Type=simple
User=@DEFAULT_USER@
Group=@DEFAULT_GROUP@
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/minidlnad -S
ProtectSystem=full
ProtectHome=on
PrivateDevices=on
NoNewPrivileges=on

[Install]
WantedBy=multi-user.target
